# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:camofox]
# Run the Camofox Browser fork directly. The compiled fork lives under
# /opt/camofox/ in the runtime image (Stage 4 installs the build output
# there); the entry point is dist/src/server.js per the fork's package.json.
command=node --max-old-space-size=128 dist/src/server.js
directory=/opt/camofox
user=node
environment=HOME="/home/node",NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:opencli-daemon]
command=node /opt/opencli/dist/src/daemon.js
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2

[program:shim]
command=node /opt/shim/index.js
user=node
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=3
depends_on=camofox opencli-daemon

[program:gateway]
command=node /opt/gateway/mcp/index.js
user=node
# CAMOFOX_API_KEY/URL/USER_ID are forwarded so the video subsystem can pull
# cookies from /sessions/<userId>/cookies (used by yt-dlp-based downloads).
# Without CAMOFOX_API_KEY the Camofox cookies endpoint returns 403 and every
# non-native download fails (LOGIN_REQUIRED even when the user is logged in).
# The %(ENV_*)s syntax pulls values from the supervisord process's own
# environment, which inherits the values docker-compose.yml injects into the
# container — keeps the secret in one place.
environment=HOME="/home/node",NODE_ENV="production",GATEWAY_LOG_DIR="/var/log/gateway",CAMOFOX_API_KEY="%(ENV_CAMOFOX_API_KEY)s",CAMOFOX_URL="%(ENV_CAMOFOX_URL)s",CAMOFOX_USER_ID="%(ENV_CAMOFOX_USER_ID)s"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/gateway/supervisor-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/gateway/supervisor-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=4
depends_on=opencli-daemon
