# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:camofox]
# Run the Camofox Browser fork directly. The Dockerfile installs the
# compiled fork at /opt/camofox (Stage 4), not /app — supervisord
# requires `directory=` to be set so node picks up the correct
# node_modules and any relative paths.
command=node --max-old-space-size=128 server.js
directory=/opt/camofox
user=node
environment=HOME="/home/node",NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:opencli-daemon]
command=node /opt/opencli/dist/src/daemon.js
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2

[program:shim]
command=node /opt/shim/dist/index.js
user=node
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=3
depends_on=camofox opencli-daemon

[program:gateway]
command=node /opt/gateway/dist/index.js
user=node
environment=HOME="/home/node",GATEWAY_LOG_DIR="/var/log/gateway"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/gateway/supervisor-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/gateway/supervisor-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=4
depends_on=opencli-daemon
